# my-cookbooks/errbit/templates/errbit.service.erb
[Unit]
Description=<%= node['errbit']['systemd']['service']['description'] %>

[Service]
Type=simple
User=<%= @user %>
WorkingDirectory=<%= @dir %>
ExecStart=/bin/bash -lc 'source ~/.profile && bundle exec rails server -p <%= @port %> -b <%= @interface %>'
Restart=always

[Install]
WantedBy=multi-user.target

